# deploy/syncrow-portal.service
[Unit]
Description=SyncRow Portal (FastAPI + HTMX + Plotly)
# InfluxDB runs on this same box; start after it and the network are up.
After=network-online.target influxdb.service
Wants=network-online.target

[Service]
Type=simple
User=syncrow
Group=syncrow
WorkingDirectory=/opt/syncrow/SyncRow_Portal

# InfluxDB creds + config (root-owned, mode 640), NOT in the repo.
# Provides: INFLUX_URL, INFLUX_TOKEN, INFLUX_ORG, INFLUX_ORG_ID, INFLUX_BUCKET.
EnvironmentFile=/etc/syncrow/portal.env

# Bind to localhost only — nginx is the sole public face (TLS + basic-auth via
# /etc/nginx/.htpasswd terminate there).
ExecStart=/opt/syncrow/SyncRow_Portal/.venv/bin/uvicorn web.main:app \
    --host 127.0.0.1 \
    --port 5006

Restart=on-failure
RestartSec=5

# --- Hardening ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ReadWritePaths=/opt/syncrow/SyncRow_Portal

[Install]
WantedBy=multi-user.target
